# templates/default/jupyter.conf.erb
description "jupyter service"
author "Sammy Lin <sammylintw@gmail.com>"
start on startup
stop on shutdown

respawn
respawn limit 10 5

script
  su - <%= node['deploy']['deploy_user'] %> -c "cd <%= @deploy_path %> && \
  . .virtualenv/bin/activate && \
  jupyter notebook --ip 0.0.0.0"
end script
